# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisor/supervisord.pid
childlogdir=/var/log/supervisor
user=root

[program:openclaw-gateway]
command=node /app/openclaw.mjs gateway --allow-unconfigured
directory=/app
user=node
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=NODE_ENV="production",HOME="/home/node",OPENCLAW_GATEWAY_PORT="%(ENV_OPENCLAW_GATEWAY_PORT)s",OPENCLAW_BRIDGE_PORT="%(ENV_OPENCLAW_BRIDGE_PORT)s",OPENCLAW_GATEWAY_TOKEN="%(ENV_OPENCLAW_GATEWAY_TOKEN)s",OPENCLAW_ALLOW_INSECURE_PRIVATE_WS="%(ENV_OPENCLAW_ALLOW_INSECURE_PRIVATE_WS)s",OPENCLAW_GATEWAY_URL="%(ENV_OPENCLAW_GATEWAY_URL)s",ZAI_API_KEY="%(ENV_ZAI_API_KEY)s",MINIMAX_API_KEY="%(ENV_MINIMAX_API_KEY)s",ANTHROPIC_AUTH_TOKEN="%(ENV_ANTHROPIC_AUTH_TOKEN)s",ANTHROPIC_BASE_URL="%(ENV_ANTHROPIC_BASE_URL)s",NOTION_API_TOKEN="%(ENV_NOTION_API_TOKEN)s",BRAVE_API_KEY="%(ENV_BRAVE_API_KEY)s",DEEPGRAM_API_KEY="%(ENV_DEEPGRAM_API_KEY)s",OPENAI_API_KEY="%(ENV_OPENAI_API_KEY)s",DISCORD_BOT_TOKEN="%(ENV_DISCORD_BOT_TOKEN)s",GOOGLE_API_KEY="%(ENV_GOOGLE_API_KEY)s",SUPABASE_PROJECT_URL="%(ENV_SUPABASE_PROJECT_URL)s",SUPABASE_PROJECT_API="%(ENV_SUPABASE_PROJECT_API)s",SUPABASE_CONNECTION_STRING_IPV4="%(ENV_SUPABASE_CONNECTION_STRING_IPV4)s",DISCORD_USER_ID="%(ENV_DISCORD_USER_ID)s",DISCORD_BOT_PERMISSIONS="%(ENV_DISCORD_BOT_PERMISSIONS)s",DISCORD_MORNING_CHANNEL_ID="%(ENV_DISCORD_MORNING_CHANNEL_ID)s",DISCORD_RAIN_ALERT_CHANNEL_ID="%(ENV_DISCORD_RAIN_ALERT_CHANNEL_ID)s",DISCORD_CHANNEL_ID="%(ENV_DISCORD_CHANNEL_ID)s",DISCORD_CHANNEL_ID_MUSIC="%(ENV_DISCORD_CHANNEL_ID_MUSIC)s",DISCORD_CHANNEL_ID_SHOPEE="%(ENV_DISCORD_CHANNEL_ID_SHOPEE)s",LINE_CHANNEL_ACCESS_TOKEN="%(ENV_LINE_CHANNEL_ACCESS_TOKEN)s",LINE_CHANNEL_SECRET="%(ENV_LINE_CHANNEL_SECRET)s",LINE_USER_PI_PRAO="%(ENV_LINE_USER_PI_PRAO)s",LINE_USER_CHUNNY="%(ENV_LINE_USER_CHUNNY)s"

; Typhoon ASR API - Speech-to-Text service
; Runs from workspace/skills/music-class-summarizer/typhoon-api.py
[program:typhoon-api]
command=python3 -m uvicorn typhoon-api:app --host 0.0.0.0 --port 8000
directory=/home/node/.openclaw/workspace/skills/music-class-summarizer
user=node
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1",TYPHOON_MODEL_PATH="/typhoon",OPENAI_API_KEY="%(ENV_OPENAI_API_KEY)s",ANTHROPIC_API_KEY="%(ENV_ANTHROPIC_API_KEY)s",GOOGLE_API_KEY="%(ENV_GOOGLE_API_KEY)s"
